# modules/install-k8s/kubelet.path
[Path]
PathExists=/etc/sysconfig/kubernetes.conf

[Install]
WantedBy=multi-user.target